# config/voicebot-bot.service
[Unit]
Description=Voicebot Pipecat Bot (Twilio Voice Pipeline)
Documentation=https://github.com/dcodingx/APP_NEC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__BOT_DIR__
EnvironmentFile=__DEPLOY_DIR__/config/client.env
ExecStart=__BOT_DIR__/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 5000 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicebot-bot

[Install]
WantedBy=multi-user.target
